# api/scripts/services/api.conf
[program:nameko]
command=/bin/bash -c /launch.sh
directory=/usr/src/app
redirect_stderr=true
