# scripts/services/nameko.conf
[program:nameko]
command=/bin/bash -c /launch.sh
directory=/usr/src/app
