# services/nameko.conf
[program:nameko]
command=/usr/bin/nameko %(ENV_ACTION)s %(ENV_SERVICE_NAME)s
toto
directory=/usr/src/app
