# services/nameko.conf
[program:nameko]
command=nameko %(ENV_ACTION)s --config config.yml %(ENV_SERVICE_NAME)s
directory=/usr/src/app
